# deployment/systemd/africgraph-backup.service
[Unit]
Description=AfricGraph Backup Service
After=network.target docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/africgraph
EnvironmentFile=/opt/africgraph/.env
ExecStart=/opt/africgraph/deployment/scripts/backup-enhanced.sh incremental
StandardOutput=journal
StandardError=journal
